# pihole-ha-sync-pull.service
[Unit]
Description=Pi-hole HA Config Sync Pull from Primary
After=network-online.target pihole-FTL.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pihole-ha-sync-pull
StandardOutput=journal
StandardError=journal
